# supervisord.conf
[supervisord]
nodaemon=true

[program:rtorrent]
priority=1
command=/root/startup-rtorrent.sh
redirect_stderr=true
